# docker/local/supervisor/consumer.conf
[program:consumer]
command=npm run consumer
directory=/usr/src/app
autostart=true
autorestart=true
stderr_logfile=/var/log/consumer.err.log
stdout_logfile=/var/log/consumer.out.log
